# piServerProfile/default.cron
#Default cron jobs

*/1 * * * * /root/tunnel.sh

