# piServerProfile/sync.cron
# Sync
*/10 * * * * /home/pi/piServerProfile/sync.sh

