# piServerProfile/tv.cron
#PiServer TV

0 7 * * 1,2,3,4,5 /home/pi/piServerProfile/tv.sh
2 8 * * 1,2,3,4,5 /home/pi/piServerProfile/tv.sh
